# templates/nsupdate.service.j2
# generated by ansible for host {{ inventory_hostname }}

[Unit]
Description=nsupdate
Requires=network.target
After=syslog.target network.target

[Service]
#Type=forking
ExecStart={{ nsupdate_installation_path }}/nsupdate-run.sh
User={{ nsupdate_user }}

[Install]
WantedBy=multi-user.target
